# phoromatic.ini
[program:phoromatic]
nodaemon=true
command=phoronix-test-suite phoromatic.connect %(ENV_SERVER_ADDRESS)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
